# roles/netbird-server/templates/netbird.service.j2
[Unit]
Description=NetBird Management Service
After=network.target

[Service]
EnvironmentFile=/etc/netbird/env
ExecStart=/usr/local/bin/netbird-mgmt
Restart=on-failure

[Install]
WantedBy=multi-user.target
